# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
stdout_syslog=true
autorestart=true
priority=10